# aws-ecs/packer/to-upload/weave.conf
# Copyright 2014-2015 Weaveworks. All Rights Reserved.
#
# Licensed under the Apache License, Version 2.0 (the
# "License"). You may not use this file except in compliance
# with the License. A copy of the License is located at
#
#     http://aws.amazon.com/apache2.0/
#
# or in the "license" file accompanying this file. This file is
# distributed on an "AS IS" BASIS, WITHOUT WARRANTIES OR
# CONDITIONS OF ANY KIND, either express or implied. See the
# License for the specific language governing permissions and
# limitations under the License.

description "Weave Service init"
author "Weaveworks"
start on started docker
stop on stopping docker

script
  while true; do
    # verify that weave is not running
    while [ "$(docker inspect -f '{{.State.Running}}' weaveproxy 2> /dev/null )" = true ]; do sleep 2; done
    # launch weave
    peers=$(bash /etc/weave/peers.sh)
    if [ ! -z "$peers" ]; then
      weave launch-router $peers
    else
      weave launch-router
    fi
    # ecs-agent uses the following derived container names when launching tasks:
    # taskFamiliyName-Version-ContainerName-UUID
    #
    # For instance: task-example-5-httpserver-aceb93e2f2b797caba01
    #
    # For now, let's use --hostname-match to force registering the original
    # ContainerName as a hostname in weaveDNS
    weave launch-proxy --hostname-match '^.+-[0-9]+-(.+)-[0-9a-f]+$'
  done
end script

pre-stop exec weave stop


